# systemd/renewIp.timer
[Unit]
Description=renew dynamic ip via cloudflare every few minutes

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=timers.target