# systemd/renewIp.service
[Unit]
Description=renew dynamic ip via cloudflare

[Service]
ExecStart=/root/nimCloudflareDynDns/renewIP
WorkingDirectory=/root/nimCloudflareDynDns/